# base/golang/supervisord.conf
[supervisord]
nodaemon=true

[program:goproxy]
command=/bin/bash -c "source ~/.profile && service ssh start"
